# systemd/docker-volume-onedata.socket
[Unit]
Description=Onedata Docker volume plugin

[Socket]
ListenStream=/run/docker/plugins/onedata.sock

[Install]
WantedBy=sockets.target
